# userpatches/overlay/etc/systemd/system/octessera-update-recovery.service
[Unit]
Description=Octessera conservative update recovery
Before=octessera.service

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/sbin/octessera-update-recovery
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
