# tools/pi-image/stage4-octessera/files/root/etc/systemd/system/octessera-sd-card.service
[Unit]
Description=octessera OLED SD card mount
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/octessera-sd-card mount

[Install]
WantedBy=multi-user.target
