# tools/pi-image/stage4-octessera/files/root/etc/systemd/system/octessera-network-health.service
[Unit]
Description=octessera Network Health Logger
After=NetworkManager.service ssh.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/octessera-network-health
