# tools/pi-image/stage4-octessera/files/root/etc/systemd/system/octessera-network-health.timer
[Unit]
Description=Run octessera network health checks

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
