# tools/pi-image/stage4-cellsymphony/files/root/etc/systemd/system/cellsymphony.service
[Unit]
Description=Cell Symphony Pi Zero 2W Headless Music System
After=sound.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
Environment=CELLSYMPHONY_PI_SAMPLES_DIR=/home/pi/samples
Environment=CELLSYMPHONY_PI_STORE_DIR=/home/pi/presets
ExecStartPre=/bin/sleep 2
ExecStart=/usr/local/bin/cellsymphony-pi
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
