# tools/pi-image/stage4-cellsymphony/files/root/etc/systemd/system/cellsymphony-boot-splash.service
[Unit]
Description=Cell Symphony Early OLED Boot Splash
DefaultDependencies=no
After=systemd-modules-load.service systemd-udevd.service
Before=sysinit.target cellsymphony.service

[Service]
Type=oneshot
ExecStart=-/usr/local/bin/cellsymphony-pi --boot-splash-once
TimeoutStartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=sysinit.target
